# config/chroot_local-includes/usr/lib/systemd/user/tails-persistence-is-enabled.target
[Unit]
Description=Tails persistence is enabled
Documentation=https://tails.boum.org/contribute/design/persistence
