# config/chroot_local-includes/lib/systemd/system/tails-persistent-storage.service
[Service]
Type=dbus
BusName=org.boum.tails.PersistentStorage
ExecStart=/usr/local/lib/tps
# XXX: Add sandboxing options
